# gui.service
[Unit]
Description=gui-service

[Service]
Type=oneshot
ExecStart=sudo bash /usr/local/bin/gui.sh

[Install]
WantedBy=multi-user.target 